# camera_restart.timer
[Unit]
Description=30 minute timer

[Timer]
OnCalendar=*:0/30
Unit=camera_restart.service

[Install]
WantedBy=basic.target  
