# camera_restart.timer
[Unit]
Description=5 minute timer

[Timer]
OnCalendar=*:0/5
Unit=camera_restart.service

[Install]
WantedBy=basic.target  
